# files/systemd/change-user-shell-zsh.service
[Unit]
Description=Change user shells to zsh (one-time per deployment)
After=local-fs.target systemd-user-sessions.service
Before=getty.target
ConditionPathExists=/usr/bin/zsh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/change-user-shell-zsh.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
