# files/systemd/fix-vpn-selinux.service
[Unit]
Description=Fix VPN certificate SELinux context
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/fix-vpn-selinux.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
