# files/systemd/fix-vpn-selinux.service
[Unit]
Description=Fix OpenVPN Certificate SELinux Context
After=NetworkManager.service
Documentation=man:chcon(1)

[Service]
Type=oneshot
EnvironmentFile=/etc/fix-vpn-selinux.env
ExecStart=/usr/bin/chcon -R -t ${VPN_CERT_CONTEXT} ${VPN_CERT_PATH}
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
