# files/systemd/fix-vpn-selinux.service
[Unit]
Description=Fix VPN certificate SELinux context
After=NetworkManager.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/fix-vpn-selinux.env
ExecStart=/usr/bin/chcon -R -t ${VPN_CERT_CONTEXT:-home_cert_t} ${VPN_CERT_PATH:-/var/home/%u/.cert}
ExecStart=/usr/bin/chcon -R -t ${VPN_CERT_CONTEXT:-home_cert_t} ${VPN_NM_CERT_PATH:-/var/home/%u/.local/share/networkmanagement/certificates/nm-openvpn}
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
