# files/systemd/fix-vpn-selinux.service
[Unit]
Description=Fix VPN certificate SELinux context (permanent)
After=local-fs.target
Before=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix-vpn-selinux.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
