# files/system/etc/systemd/system/flatpak-auto-update.service
[Unit]
Description=Weekly Flatpak auto-update
Documentation=man:flatpak(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak-auto-update.sh
