# files/system/etc/systemd/system/install-surface-flatpaks.service
[Unit]
Description=Install Surface Flatpaks after deployment
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/install-surface-flatpaks.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
