# files/system/etc/systemd/system/install-common-flatpaks.service
[Unit]
Description=Install common Flatpak applications (one-time per deployment)
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/justin-os/common-flatpaks-installed

[Service]
Type=oneshot
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
ExecStart=/usr/bin/install-common-flatpaks.sh

[Install]
WantedBy=multi-user.target
